# config/room_guard.service
[Unit]
Description=Room Guard — Motion Detection Alarm
After=multi-user.target

[Service]
Type=simple
User=yehudalevavi
WorkingDirectory=/home/yehudalevavi/rpiProject
ExecStart=/home/yehudalevavi/rpiProject/.venv/bin/python3 /home/yehudalevavi/rpiProject/src/room_guard.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
